# com.expanse.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.expanse.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/expanse_daemon.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/expanse_daemon.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/expanse_daemon.out</string>
</dict>
</plist> 